# dockerfiles/packaging/debian/james.postinst
#! /bin/sh -e
 
useradd --home-dir /var/lib/james james

chown -R james:james /var/log/james
chown -R james:james /var/lib/james

ln -s /etc/james /var/lib/james/conf

systemctl enable james

