# dockerfiles/packaging/guice/cassandra/scripts/james.rpm.postinst
#! /bin/sh -e
 
ln -s /etc/james /var/lib/james/conf

#Workaround waiting for https://github.com/jordansissel/fpm/issues/1163 to be released
cp /usr/share/james/james.service /etc/systemd/system/
systemctl enable james

