# dockerfiles/packaging/guice/cassandra/scripts/james.postinst
#! /bin/sh -e

ln -s /etc/james /var/lib/james/conf

systemctl enable james

